# rf-cdp/services/ingest-gateway-prod/cdp-gateway.service
[Unit]
Description=CDP ingest-gateway
After=network.target
[Service]
Type=simple
WorkingDirectory=/opt/cdp-gateway
EnvironmentFile=/opt/cdp-gateway/.env
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=3
User=root
[Install]
WantedBy=multi-user.target
